# paperclipai.service
[Unit]
Description=PaperclipAI Local Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aa
Group=aa
WorkingDirectory=/home/aa/vllm
ExecStart=/home/aa/.npm-global/bin/paperclipai run --data-dir ./paperclip-data
Restart=always
RestartSec=5
Environment=PATH=/home/aa/.local/bin:/home/aa/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/aa
Environment=NODE_ENV=production

# 휴지(sleep/hibernate) 모드 복귀 후에도 자동 재시작
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
